# scripts/com.quant.daily_job_runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 1. 唯一标识符，必须与文件名对应 -->
    <key>Label</key>
    <string>com.yourcompany.daily-job-runner</string>

    <!-- 2. 要执行的程序和参数 -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/你的用户名/miniconda3/bin/python3</string> <!-- 替换为你的 Python 路径 -->
        <string>/Users/你的用户名/项目路径/backend/scripts/daily_job_runner.py</string> <!-- 替换为脚本绝对路径 -->
        <string>--stage</string>
        <string>1</string>
        <string>--fill-missing</string>
    </array>

    <!-- 3. 工作目录（非常重要！确保脚本中的相对路径如 logs/ 正确） -->
    <key>WorkingDirectory</key>
    <string>/Users/你的用户名/项目路径</string>

    <!-- 4. 环境变量（让脚本能找到 .env 和第三方库） -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/你的用户名/miniconda3/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/你的用户名/项目路径</string>
    </dict>

    <!-- 5. 日志输出（方便排查问题） -->
    <key>StandardOutPath</key>
    <string>/Users/你的用户名/项目路径/logs/cron/runner.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/你的用户名/项目路径/logs/cron/runner.stderr.log</string>

    <!-- 6. 定时触发：每周一到周五的 15:30 和 16:30 执行 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer> <!-- 周一 -->
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- 周二到周五按相同格式重复 -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <!-- 7. 防止任务重叠（如果前一个没跑完，新的不启动） -->
    <key>ThrottleInterval</key>
    <integer>0</integer>

    <!-- 8. 失败时最多尝试重启次数 -->
    <key>NumberOfRetries</key>
    <integer>3</integer>
</dict>
</plist>