# backend/scheduler/launchd/daily_pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quant.daily_pipeline</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zhangk/workspace/Quantitative_trading/backend/scheduler/run_daily_pipeline.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/zhangk/workspace/Quantitative_trading</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/zhangk/workspace/Quantitative_trading/logs/scheduler/daily_pipeline.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhangk/workspace/Quantitative_trading/logs/scheduler/daily_pipeline.err.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>30</integer>
        <key>Weekday</key>
        <array>
            <integer>1</integer>
            <integer>2</integer>
            <integer>3</integer>
            <integer>4</integer>
            <integer>5</integer>
        </array>
    </dict>
</dict>
</plist>